# puppet/modules/go-ethereum/templates/go-ethereum.upstart.conf.erb
description "go-ethereum service"

start on runlevel [23]
stop on shutdown

setuid ethereum
setgid ethereum

kill signal INT

exec <%= @daemon_path %> -mine -maxpeer <%= @max_peer %> -port <%= @outbound_port %> -logfile <%= @log_file %>

respawn